# window-shadow.service
[Unit]
Description=X11 Window Shadow Daemon
After=x11-common.service graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/window-shadow
Restart=on-failure
RestartSec=2
# In a user session, DISPLAY and XAUTHORITY are usually already set,
# but we provide a fallback just in case.
Environment=DISPLAY=:0

[Install]
WantedBy=graphical-session.target
