# window-shadow.service
[Unit]
Description=X11 Window Shadow Daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/window-shadow
Restart=always
RestartSec=5
# In a user session, DISPLAY and XAUTHORITY are usually already set,
# but we provide a fallback just in case.
Environment=DISPLAY=:0

[Install]
WantedBy=default.target
